# systemd/vinput-daemon.service.in
[Unit]
Description=Vinput Voice Input Daemon
After=pipewire.service

[Service]
Type=dbus
BusName=org.fcitx.Vinput
ExecStart=@VINPUT_DAEMON_EXECUTABLE@

[Install]
WantedBy=default.target
